# deploy/discord-bot.service
# Copy to server: sudo cp deploy/yuuki.service /etc/systemd/system/
# Then: sudo systemctl daemon-reload && sudo systemctl enable discord-bot
# Edit paths and user below to match your server.

[Unit]
Description=Discord bot
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/projects/app
Environment=PATH=/home/projects/app.venv/bin:/usr/bin
ExecStart=/home/projects/app/.venv/bin/python bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
